# _harp/crontab
30 20 * * * python /Users/Marek/Downloads/dribbble.py
